# crontab
*/1 * * * * /usr/local/bin/php /var/www/html/cron/status.cron.php --timeout=310
